# adapters/claude-code/ai.warden.gateway.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Template. install.sh renders this to ~/Library/LaunchAgents/ai.warden.gateway.plist,
     substituting __WARDEN_BIN__ (the installed warden binary) and __HOME__ (your
     absolute home dir — launchd does NOT expand $HOME). Do not load this file
     directly; run install.sh. -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.warden.gateway</string>
    <key>ProgramArguments</key>
    <array>
        <string>__WARDEN_BIN__</string>
        <string>serve</string>
        <string>--socket</string>
        <string>/tmp/warden_daemon.sock</string>
        <string>--policy</string>
        <string>__HOME__/.warden/policy.tsv</string>
        <string>--manifests</string>
        <string>__HOME__/.warden/manifests</string>
        <string>--audit</string>
        <string>__HOME__/.warden/audit.jsonl</string>
        <string>--non-interactive</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardErrorPath</key>
    <string>__HOME__/.warden/gateway.err.log</string>
    <key>StandardOutPath</key>
    <string>__HOME__/.warden/gateway.out.log</string>
</dict>
</plist>
